# install/systemd/claude-signal.service
[Unit]
Description=Claude Signal - Signal to Claude Code Bridge
Documentation=https://github.com/dmz006/claude-signal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude-signal
Group=claude-signal
WorkingDirectory=/var/lib/claude-signal
ExecStart=/usr/local/bin/claude-signal start --config /etc/claude-signal/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-signal

NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=/var/lib/claude-signal /var/log/claude-signal /tmp
AmbientCapabilities=

Environment=HOME=/var/lib/claude-signal
Environment=XDG_DATA_HOME=/var/lib/claude-signal/.local/share

[Install]
WantedBy=multi-user.target
